# app/web/supervisord.conf
[supervisord]
nodaemon=true

[program:shell]
command=/bin/bash

[program:sshd]
command=/usr/sbin/sshd -D


